# supervisor/deluged.conf
[program:deluged]
command=/usr/bin/deluged -u 127.0.0.1 -d
autostart=true