# supervisor/deluge.conf
[program:deluge]
command=/usr/bin/deluge-web -i 127.0.0.1 -p 8112 -d
autostart=true