# supervisor/filebrowser_noauth.conf
[program:filebrowser_noauth]
directory=/home/ubuntu/filebrowser/noauth
command=/usr/local/bin/filebrowser --noauth
user=ubuntu
autostart=true
autorestart=true
