# supervisor/filebrowser.conf
[program:filebrowser]
directory=/home/ubuntu/filebrowser
command=/usr/local/bin/filebrowser
user=ubuntu
autostart=true
autorestart=true
