# supervisor/filebrowser.conf
[program:filebrowser]
directory=/home/ubuntu/filebrowser
command=/usr/local/bin/filebrowser
autostart=true
autorestart=true
