# supervisor/frp.conf
[program:frp]
directory=/home/ubuntu/frp
command=/home/ubuntu/frp/frps -c /home/ubuntu/frp/frps.ini
autostart=true
