# supervisor/glances.conf
[program:glances]
command=/usr/local/bin/glances -w
user=ubuntu
